# etc/systemd/earlyoom.service.d/memory-pressure.conf
[Service]
ExecStart=
ExecStart=/usr/bin/earlyoom -r 3600 -n -m 10,5 -s 30,20 --avoid (^|/)(init|systemd|sshd|Hyprland|xdg-desktop-portal-hyprland|eww|eww-bar)$$
