# config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session daemons
Wants=hyprd.service ewwd.service
