# config/systemd/user/opencode.service
[Unit]
Description=OpenCode shared backend

[Service]
ExecStart=/usr/bin/opencode serve --hostname 127.0.0.1 --port 4096
Restart=on-failure
RestartSec=2
Nice=10
CPUWeight=25
IOWeight=25
CPUQuota=1000%
