# config/systemd/user/pipewire.service.d/cpu-lane.conf
[Service]
CPUAffinity=7 15
CPUWeight=1000
